# src/crontab
0 */3 * * * /path/to/php path/to/Sessiongc.php